# kasm/rootfs/etc/s6-overlay/s6-rc.d/init-hassio/run
#!/usr/bin/bashio

mkdir -p /config/profiles
rm -r /profiles
ln -s /config/profiles /

mkdir -p /data/opt
rm -r /opt
ln -s /data/opt /