# tubesync/rootfs/etc/s6-overlay/s6-rc.d/hassio/run
#!/command/with-contenv /usr/bin/bashio

if bashio::config.has_value 'TZ'; then
    bashio::config 'TZ' > /var/run/s6/container_environment/TZ


fi

# if bashio::config.has_value 'HTTP_USER'; then
#     (bashio::config 'HTTP_USER') > /var/run/s6/container_environment/HTTP_USER
#     (bashio::config 'HTTP_PASS') > /var/run/s6/container_environment/HTTP_PASS
# fi

dlf="$(bashio::config 'downloadFolder')"

mkdir -p "$dlf"/video
mkdir -p "$dlf"/audio
rm -rf /downloads/video
rm -rf /downloads/audio
chown -R 911 "$dlf"/*
ln -s "$dlf"/* /downloads/