# tubesync/rootfs/etc/s6-overlay/s6-rc.d/hassio/run
#!/command/with-contenv /usr/bin/bashio

if bashio::config.has_value 'TZ'; then
    export TZ=$(bashio::config 'TZ')
fi

if bashio::config.has_value 'HTTP_USER'; then
    export HTTP_USER=$(bashio::config 'HTTP_USER')
    export HTTP_PASS=$(bashio::config 'HTTP_PASS')
fi

mkdir -p $(bashio::config 'downloadFolder')/video
mkdir -p $(bashio::config 'downloadFolder')/audio
rm -rf /downloads/video
rm -rf /downloads/audio
ln -s $(bashio::config 'downloadFolder')/* /downloads/